# com.madewell.autopush.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.madewell.autopush</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/madewell/Documents/Claude/Projects/MADEWELL MUSIC 리브랜딩/autopush.sh</string>
  </array>
  <key>StartInterval</key>
  <integer>1800</integer>
  <key>RunAtLoad</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/madewell/Documents/Claude/Projects/MADEWELL MUSIC 리브랜딩/.autopush.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/madewell/Documents/Claude/Projects/MADEWELL MUSIC 리브랜딩/.autopush.log</string>
</dict>
</plist>
